# cookbooks/atlassian/templates/default/bamboo.erb
#!/bin/bash

#
# bamboo          Start up the Bamboo server
#
# chkconfig: 2345 55 25
# description: Bamboo is an automated build server
#
# processname: bamboo

# Bamboo Linux service controller script
cd "/opt/atlassian/bamboo"

case "$1" in
    start)
        su -c '/opt/atlassian/bamboo/bamboo.sh start' - bamboo
        ;;
    stop)
        su -c '/opt/atlassian/bamboo/bamboo.sh stop' - bamboo
        ;;
    status)
        su -c '/opt/atlassian/bamboo/bamboo.sh status' - bamboo
        ;;
    restart)
        su -c '/opt/atlassian/bamboo/bamboo.sh restart' - bamboo
        ;;
    *)
        echo "Usage: $0 {start|stop|status|restart}"
        exit 1
        ;;
esac

