# cookbooks/mongodb/templates/default/mongod-init.erb
#!/bin/bash

# <%= @app_name %> - Startup script for mongod

# chkconfig: 35 85 15
# description: Mongo is a scalable, document-oriented database.
# processname: mongod
# config: /etc/<%= @app_name %>.conf
# pidfile: /var/run/mongo/<%= @app_name %>.pid

. /etc/rc.d/init.d/functions

# things from mongod.conf get there by mongod reading it


# NOTE: if you change any OPTIONS here, you get what you pay for:
# this script assumes all options are in the config file.
CONFIGFILE="/etc/<%= @app_name %>.conf"
PIDFILE="/var/run/mongo/<%= @app_name %>.pid"
LOCKFILE="/var/lock/subsys/<%= @app_name %>"
OPTIONS=" -f $CONFIGFILE --pidfilepath=$PIDFILE"
SYSCONFIG="/etc/sysconfig/mongod"
mongod=/usr/bin/mongod
#mongod=${MONGOD-/usr/bin/mongod}
prog=<%= @app_name %>

MONGO_USER=mongod
MONGO_GROUP=mongod

. "$SYSCONFIG" || true

start()
{
  echo -n $"Starting $prog: "
  daemon --user "$MONGO_USER" $mongod $OPTIONS
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && touch $LOCKFILE
}

stop()
{
  echo -n $"Stopping $prog: "
  killproc -p $PIDFILE -d 30 -TERM $mongod
  RETVAL=$?
  echo
  [ $RETVAL -eq 0 ] && rm -f $LOCKFILE
}

restart () {
	stop; sleep 5
	start
}

ulimit -n 12000
RETVAL=0

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|reload|force-reload)
    restart
    ;;
  condrestart)
    [ -f $LOCKFILE ] && restart || :
    ;;
  status)
    status -p $PIDFILE $mongod
    RETVAL=$?
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart|reload|force-reload|condrestart}"
    RETVAL=1
esac

exit $RETVAL
