# cookbooks/infrastructure/templates/default/appdynamics-init.erb
#!/bin/bash

# appdynamics          Start up the Appdynamics Controller daemon
#
# chkconfig: 2345 90 10
# description: Appdynamics Controller. \
#
# processname: appdynamics
# pidfile: /var/run/appdynamics.pid

# Appdynamics Controller Linux service controller script
cd "/opt/appdynamics/bin"

case "$1" in
    start)
        su -c "/opt/appdynamics/bin/controller.sh start" appdynamics
        ;;
    stop)
        su -c "/opt/appdynamics/bin/controller.sh stop" appdynamics
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    status)
        pid_status=`ps -ef | grep java | grep appdynamics | awk '{print $2}'`
        if [ "$pid_status" != "" ]; then
          echo "Appdynamics Controller is running with pid $pid_status."
        else
          echo "Appdynamics Controller is not running."
        fi
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

