# cookbooks/altisource/templates/default/splunkforwarder.erb
#!/bin/bash

# splunk          Start up the Splunk server daemon
#
# chkconfig: 2345 55 25
# description: Splunk is a server for log file analysis. \
#              This service starts up the Splunk server daemon.
#
# processname: python
# pidfile: /var/run/splunk.pid

# Splunk Linux service controller script
SPLUNK_HOME=/opt/splunkforwarder; export SPLUNK_HOME

cd "/opt/splunkforwarder/bin"

case "$1" in
    start)
        su -c "nohup ./splunk start &" splunk > /dev/null 2>&1
        ;;
    stop)
        su -c "./splunk stop" splunk
        ;;
    restart)
	$0 stop
	$0 start
        ;;
    status)
        pid_status=`ps -ef | grep splunkd | grep -v grep | awk '{print $2}'`
        if [ "$pid_status" != "" ]; then
          echo "Splunk is running with pid $pid_status."
        else
          echo "Splunk is not running."
        fi
        ;;
    *) 
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

