# cookbooks/realdoc/templates/default/realdoc-correspondence-init.erb
#!/bin/sh
#
# `basename $0`
#
# chkconfig: - 85 15
# description: Realdoc-correspondence Java App
# processname: java
# pidfile: /var/run/realdoc-correspondence.pid

# Realdoc Correspondence Linux service controller script
. /etc/rc.d/init.d/functions
USER=javaapp
CMD=`echo ${0##*/}`
LOCKFILE=/var/lock/realdoc-correspondence
PIDFILE=/var/run/realdoc-correspondence.pid

case "$1" in
    start)
        echo -n "Starting Realdoc-correspondence: "
	status -p $PIDFILE $CMD > /dev/null && failure || (su -p -s /bin/sh $USER -c "/opt/realdoc-correspondence/bin/realdoc-correspondence" > /dev/null && (touch $LOCKFILE ; success))
	echo
	;;
    stop)
	echo -n "Shutting down $CMD: "
	PID=`cat $PIDFILE`; kill $PID
	SLEEP=20
	while [ $SLEEP -gt 0 ]; do
	  stop_status=`ps -fu $USER|fgrep java|awk '{print $2}'`
	  if [ "$stop_status" != "" ]; then
            echo "."
          elif [ "$stop_status" == "" ]; then
            success; break
	  fi
	  SLEEP=$(( $SLEEP - 1 ))
	done
        stop_status=`ps -fu $USER| fgrep java|awk '{print $2}'`
        if [ "$stop_status" != "" ]; then
          echo "Attempting to hard kill."; ps -fu $USER| fgrep java|awk '{print $2}' | xargs -r kill -9 && rm -f $LOCKFILE && rm -f $PIDFILE; success
	fi
        stop_status=`ps -fu $USER| fgrep java|awk '{print $2}'`
        if [ "$stop_status" == "" ]; then
                success
        else
                failure
        fi
        echo 
        ;;
    restart)
        $0 stop
	$0 start
        ;;
    status)
        status -p $PIDFILE $CMD
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac
